# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:epiviz]
command=python3 /var/www/api/epiviz.py

[program:apache]
command=apachectl -D FOREGROUND